# openrc/freellmapi
#!/sbin/openrc-run
# Portable FreeLLMAPI init script.
# Override FREEMAPI_USER in /etc/conf.d/freellmapi for a different user.

description="FreeLLMAPI — Unified free LLM router"

# Auto-detect: use the first non-system user with a home directory
if [ -z "${FREEMAPI_USER:-}" ]; then
    for _u in $(getent passwd | awk -F: '$3 >= 1000 && $3 < 65534 {print $1}'); do
        if [ -d "/home/$_u" ]; then FREEMAPI_USER="$_u"; break; fi
    done
fi
: "${FREEMAPI_USER:=root}"
FREEMAPI_USER_HOME=$(getent passwd "$FREEMAPI_USER" | cut -d: -f6)
FREEMAPI_DIR="${FREEMAPI_DIR:-${FREEMAPI_USER_HOME}/repos/freellmapi}"

depend() {
    need docker net
}

start() {
    ebegin "Starting FreeLLMAPI"
    cd "$FREEMAPI_DIR" && docker compose up -d
    eend $?
}

stop() {
    ebegin "Stopping FreeLLMAPI"
    cd "$FREEMAPI_DIR" && docker compose down
    eend $?
}

status() {
    if docker ps --format '{{.Names}}' | grep -q freellmapi; then
        einfo "FreeLLMAPI is running"
        return 0
    else
        eerror "FreeLLMAPI is not running"
        return 1
    fi
}
