# openrc/ollama
#!/sbin/openrc-run
# Portable Ollama init script.
# Override OLLAMA_USER in /etc/conf.d/ollama for a different user.

description="Ollama — Local LLM server"

# Auto-detect: use the first non-system user with a home directory
if [ -z "${OLLAMA_USER:-}" ]; then
    for _u in $(getent passwd | awk -F: '$3 >= 1000 && $3 < 65534 {print $1}'); do
        if [ -d "/home/$_u" ]; then OLLAMA_USER="$_u"; break; fi
    done
fi
: "${OLLAMA_USER:=root}"
pidfile="/run/ollama.pid"
command="/usr/bin/ollama"
command_args="serve"
command_user="$OLLAMA_USER"

depend() {
    need net
}

start() {
    ebegin "Starting Ollama"
    start-stop-daemon --start --background \
        --pidfile "$pidfile" \
        --make-pidfile \
        --user "$command_user" \
        --exec "$command" \
        -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping Ollama"
    start-stop-daemon --stop --pidfile "$pidfile"
    eend $?
}
